# supervisor_conf/job.ini
[program:job]
command=sh -c "exec /root/backup_files -d \"${BACKUP_DIRS}\" -o \"${DEST_HOST}\" -p \"${DEST_PATH}\" -c \"${DEST_CRED}\" -e \"${ENC_PASS}\" -n \"${FILE_PREFIX}\" -m ${DOM} -w ${DOW}"
environment=BACKUP_DIRS="${BACKUP_DIRS}",DEST_HOST="${DEST_HOST}",DEST_PATH="${DEST_PATH}",DEST_CRED="${DEST_CRED}",SSH_OPTS="${SSH_OPTS}",ENC_PASS="${ENC_PASS}",FILE_PREFIX="${FILE_PREFIX}",DOW="${DOW}",DOM="${DOM}",DO_DOW=1,DO_DOM=1
autostart=false
autorestart=false
startsecs=0
startretries=0

